# midirouter.service
[Unit]
Description=MIDI Router

[Service]
ExecStart=/home/alarky/midirouter/main.py
Restart=always

[Install]
WantedBy=multi-user.target
